# test/tpm2-simulator/tpm2-simulator.conf
# Copyright 2014 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.
# Rewritten by Hannah for use in CrosVM.

description     "Chromium OS TPM2.0 daemon for Trunksd"
author          "chromium-os-dev@chromium.org"

start on started boot-services and stopped cr50-result and started dbus
stop on hwsec-stop-low-level-tpm-daemon-signal
respawn


oom score -100

env TPM_DEV=/dev/tpm0

pre-start script
	STATEFUL=/mnt/stateful_partition
	SIM_DIR="${STATEFUL}/unencrypted/tpm2-simulator"

	if [ ! -d "$SIM_DIR" ]; then
		logger -t tpm2-simulator "Creating simulator directory."
		mkdir -p "$SIM_DIR"
	fi

	logger -t tpm2-simulator "Moving to simulator directory."
	cd "$SIM_DIR"

	#logger -t tpm2-simulator "Cleaning up from previous runs."
	#tpm2-simulator --socker "$TPM_DEV" --cleanup
end script

expect fork

exec tpm2-simulator --socket "$TPM_DEV"