# test/tpm2-simulator/trunksd.conf
# Copyright 2014 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Chromium OS trunks daemon for TPM2.0"
author          "chromium-os-dev@chromium.org"

# TODO(b/265866896): provider better abstraction by adding trunks-pre-init
# event instead
start on started boot-services and stopped cr50-result and started dbus
stop on hwsec-stop-low-level-tpm-daemon-signal
respawn


oom score -100

# These enviroment variable may be modified in the ebuild file.
# The runtime TPM selection feature would be enabled if this variable is true.
env TPM_DYNAMIC=false

# KEY_EVICTION value is set from the ebuild file if the USE flag is enabled.
env KEY_EVICTION=""

env TRUNKSD_FREEZER_CGROUP_DIR=/sys/fs/cgroup/freezer/trunks

pre-start script
    logger -t trunksd "Starting TPM2 Simulator."
    start tpm2-simulator

    logger -t trunksd "Initializing TPM2.0"
    tpmc startup
end script

expect fork

exec trunksd

post-start script
  if [ "${KEY_EVICTION}" = true ]; then
    /usr/bin/systemd-tmpfiles --create --remove --clean \
      /usr/lib/tmpfiles.d/on-demand/trunks_freezer.conf
    echo $(status | cut -f 4 -d ' ') > \
      "${TRUNKSD_FREEZER_CGROUP_DIR}/cgroup.procs"
  fi
end script
